# backend/systemd/apollo-miner-iii.service
[Unit]
Description=apollo-miner-iii
# Runs only when the Apollo III sentinel is present (written on boot by
# backend/utils/set_UI_mode.sh when the internal hardware is detected).
ConditionPathExists=/etc/apolloapi/apollo-iii
After=network.target rc-local.service

[Service]
Type=forking
User=root
# TBD: John — confirm start/stop scripts and binary path for the Apollo III firmware.
ExecStart=/opt/apolloapi/backend/apollo-miner-iii/miner_start.sh
WorkingDirectory=/opt/apolloapi/backend/apollo-miner-iii
ExecStop=/opt/apolloapi/backend/apollo-miner-iii/miner_stop.sh
Restart=always
RestartSec=30s

[Install]
WantedBy=multi-user.target
