# backend/systemd/node.service
[Unit]
Description=Bitcoin node
Requires=apollo-bootstrap.service
Wants=ckpool.service
After=network.target rc-local.service dev-nvme0n1.device tor.service apollo-bootstrap.service
Before=ckpool.service

[Service]
Type=simple
User=futurebit
Group=futurebit
WorkingDirectory=/opt/apolloapi
EnvironmentFile=-/opt/apolloapi/.env
ExecCondition=/opt/apolloapi/backend/node/node_storage_ready.sh
ExecStart=/opt/apolloapi/backend/node/node_start.sh
ExecStop=/opt/apolloapi/backend/node/node_stop.sh
KillSignal=SIGTERM
# node_stop waits up to 600s for dbcache flushes; keep systemd outside that window.
TimeoutStopSec=700
Restart=always
RestartSec=30s

[Install]
WantedBy=multi-user.target
